# launchd/com.xavier.manager-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.xavier.manager-weekly</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/ABSOLUTE/PATH/steel_reports_automation/scripts/manager_weekly.py</string>
      <string>--input</string>
      <string>/ABSOLUTE/PATH/steel_reports_automation/config/manager_weekly_run.json</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/steel_reports_automation</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key><integer>6</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>45</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/manager-weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/manager-weekly.err</string>
    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
